# dockerfiles/kalan-web2py/supervisor-app.conf
[program:chownweb2py]
command = chown -R www-data /home/www-data/web2py

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/web2py.ini

[program:nginx-app]
command = /usr/sbin/nginx
